# launchd/dev.fanctl.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.fanctl.watch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/fanctl</string>
        <string>watch</string>
        <string>--preset</string>
        <string>cool</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/var/log/fanctl.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/fanctl.log</string>
</dict>
</plist>
